# scripts/people-worker.service
[Unit]
Description=People Analytics Worker
After=network.target

[Service]
Type=simple
WorkingDirectory=/opt/people-analytics
ExecStart=/opt/people-analytics/.venv/bin/python -m apps.worker.worker
Restart=always

[Install]
WantedBy=multi-user.target
